# deployment/coincalltrader-deribit.service
# ===========================================================================
# CoincallTrader (Deribit) — systemd Service File
#
# Runs the Deribit trading bot as a separate service alongside the
# Coincall instance.  Dashboard on port 8081 (Coincall uses 8080).
#
# Installed to: /etc/systemd/system/coincalltrader-deribit.service
# Managed via:  systemctl start|stop|restart|status coincalltrader-deribit
# Logs via:     journalctl -u coincalltrader-deribit -f
# ===========================================================================

[Unit]
Description=CoincallTrader — Deribit Options Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Application directory and command
WorkingDirectory=/opt/coincalltrader-deribit
ExecStart=/opt/coincalltrader-deribit/.venv/bin/python main.py

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=DASHBOARD_PORT=8081

# Restart policy: restart on crash, but not on clean exit (Ctrl+C / SIGTERM).
# 10s delay prevents tight crash loops from hammering the exchange API.
Restart=on-failure
RestartSec=10

# Logging: stdout/stderr go to journald automatically.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coincalltrader-deribit

# Resource limits (safety nets)
LimitNOFILE=65535
MemoryMax=1G

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/coincalltrader-deribit/logs
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
